# recipes-devtools/android-tools/android-tools/android-tools-adbd.service
[Unit]
Description=Android Debug Bridge

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/adbd
StandardOutput=null
